# deploy/pi/ansible/roles/kioskkit/templates/kioskkit.service.j2
# Managed by Ansible — do not edit manually.
[Unit]
Description=KioskKit App
After=network.target kioskkit-privileged.service
Requires=kioskkit-privileged.service

[Service]
Type=simple
User={{ kioskkit_user }}
Group={{ kioskkit_user }}
WorkingDirectory={{ kioskkit_install_dir }}/current
ExecStart=/usr/bin/node packages/kiosk-server/dist/index.js
EnvironmentFile=-{{ kioskkit_install_dir }}/.env
{% if kioskkit_device_id %}
Environment=DEVICE_ID={{ kioskkit_device_id }}
{% endif %}
Restart=always
RestartSec=5

# Sandboxing
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ kioskkit_install_dir }}/data /data/kioskkit /data/app-update /data/ota
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
