# deploy/pi/ansible/roles/kioskkit/templates/kioskkit.service.j2
# Managed by Ansible — do not edit manually.
[Unit]
Description=KioskKit App
After=network.target

[Service]
Type=simple
User={{ kioskkit_user }}
Group={{ kioskkit_user }}
WorkingDirectory={{ kioskkit_install_dir }}
ExecStart=/usr/bin/node packages/kiosk-server/dist/index.js
EnvironmentFile=-{{ kioskkit_install_dir }}/.env
{% if kioskkit_device_id %}
Environment=DEVICE_ID={{ kioskkit_device_id }}
{% endif %}
{% if kioskkit_web_server_url %}
Environment=WEB_SERVER_URL={{ kioskkit_web_server_url }}
{% endif %}
Restart=always
RestartSec=5

# Sandboxing
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ kioskkit_install_dir }}/data /data/kioskkit
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
