# system/services/kioskkit-deploy.timer
[Unit]
Description=Daily deploy check for KioskKit

[Timer]
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
