# system/services/kioskkit.service
[Unit]
Description=KioskKit App
After=network.target

[Service]
Type=simple
User=kiosk
Group=kiosk
WorkingDirectory=/opt/kioskkit
ExecStart=/usr/bin/node packages/kiosk-server/dist/index.js
EnvironmentFile=/opt/kioskkit/.env
Restart=always
RestartSec=5

# Sandboxing
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/kioskkit/data
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
