# system/services/zahumny-kiosk-deploy.service
[Unit]
Description=Zahumny Kiosk Deploy
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/zahumny-kiosk-repo/system/deploy.sh
